# com.itnews.sender.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.itnews.sender</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/turtlesoup0-macmini/Projects/itnews_sender/scripts/launchd_runner.sh</string>
        <string>--mode</string>
        <string>opr</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <!--
      [S5] StandardOutPath / StandardErrorPath 는 제거되었다.
      launchd_runner.sh 가 실행 초반에 stdout/stderr 를 월별 파일
      (logs/launchd_YYYYMM.log) 로 직접 리디렉션하고, 90일 초과 파일을
      자동 prune 한다. 이전 launchd_stderr.log / launchd_stdout.log 는
      legacy 로 남겨두고 필요 시 수동 삭제 가능.
    -->
    <key>WorkingDirectory</key>
    <string>/Users/turtlesoup0-macmini/Projects/itnews_sender</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Asia/Seoul</string>
    </dict>
</dict>
</plist>
